# linux/crontab
* * * * * python3 /media1/compressnbackup.py